# systemd/haikubot.service
[Unit]
Description=HaikuBot IRC Haiku Generator
After=network.target

[Service]
Type=simple
User=darren
Group=darren
WorkingDirectory=/opt/Scripts/HaikuBot-Python
Environment="PATH=/opt/Scripts/HaikuBot-Python/venv/bin"
ExecStart=/opt/Scripts/HaikuBot-Python/venv/bin/python -m backend.main --config /opt/Scripts/HaikuBot-Python/config.yaml
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/Scripts/HaikuBot-Python/haikubot.log
StandardError=append:/opt/Scripts/HaikuBot-Python/haikubot-error.log

[Install]
WantedBy=multi-user.target

